# check-cc-matrix-updates.service
[Unit]
Description=Check for cc_matrix_channel upstream updates
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=PATH=%h/.cargo/bin:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=NOTIFY_AGENT=your-agent
Environment=CLAUDE_CHANNEL_PLUGIN=server:matrix
Environment=MATRIX_STATE_DIR=%h/.claude/channels/matrix-your-agent
EnvironmentFile=-%h/.claude/channels/matrix-your-agent/.env
ExecStart=%h/.local/bin/check-cc-matrix-updates
TimeoutStartSec=3600
StandardOutput=journal
StandardError=journal
