# claude-agent@.service
[Unit]
Description=Claude agent - %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/jesung/GitHub/cc-agents/%i
Environment=PATH=/home/jesung/.bun/bin:/home/jesung/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=TELEGRAM_STATE_DIR=/home/jesung/.claude/channels/telegram-%i
EnvironmentFile=/home/jesung/.claude/channels/telegram-%i/.env
ExecStart=/home/jesung/bin/claude-tmux-launch %i
ExecStartPost=/bin/bash -c '/home/jesung/bin/telegram-notify %i "Back online ($CLAUDE_MODEL)."'
ExecStop=-/usr/bin/tmux -L %i kill-server
ExecStopPost=/home/jesung/bin/telegram-notify %i "Went offline — restarting in 10s."
KillMode=control-group
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
