# claude-agent@.service
# claude-agent@.service — systemd user service template for persistent Claude agents.
#
# Usage:
#   systemctl --user enable --now claude-agent@<agent-name>
#
# Customize:
#   - WorkingDirectory: path to your agent directories (one subdir per agent)
#   - Environment=PATH: adjust to include your bun and claude install locations
#
# State dir layout expected at ~/.claude/channels/telegram-<name>/:
#   .env          — contains TELEGRAM_BOT_TOKEN=<token>
#   access.json   — access control (dmPolicy, allowFrom, etc.)
#   approved/     — created automatically

[Unit]
Description=Claude agent - %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# Directory containing your agent's CLAUDE.md / instructions.
# Structure: WorkingDirectory/<agent-name>/CLAUDE.md
WorkingDirectory=%h/agents/%i

# Adjust PATH to include wherever you installed bun and claude.
Environment=PATH=%h/.bun/bin:%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Per-agent Telegram state directory and bot token.
Environment=TELEGRAM_STATE_DIR=%h/.claude/channels/telegram-%i
EnvironmentFile=%h/.claude/channels/telegram-%i/.env

# Launch claude inside a tmux session for a real PTY.
# Assumes claude-tmux-launch is installed at ~/bin/claude-tmux-launch.
ExecStart=%h/bin/claude-tmux-launch %i

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
